# compact_measurements.service
[Unit]
Description=Compact old HiveMeasurement rows into 30-min averages
After=network.target postgresql.service

[Service]
Type=oneshot
User=main
Group=www-data
WorkingDirectory=/home/main/Hack-TUES-12-Zabravih/py/website
ExecStart=/home/main/Hack-TUES-12-Zabravih/py/env/bin/python manage.py compact_measurements
StandardOutput=journal
StandardError=journal
